# scripts/com.duckbot.memory-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  duckbot-secret-scan: allowlist-file
  com.duckbot.memory-watcher.plist — TEMPLATE.

  Install via:
    ./scripts/install-macos.sh
  which will replace every __REPO_ROOT__ with the absolute path of the
  repo before copying to ~/Library/LaunchAgents/.

  Manual install (if you must):
    sed "s|__REPO_ROOT__|$REPO_ROOT|g" \
      scripts/com.duckbot.memory-watcher.plist \
      > ~/Library/LaunchAgents/com.duckbot.memory-watcher.plist
    launchctl load -w ~/Library/LaunchAgents/com.duckbot.memory-watcher.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.duckbot.memory-watcher</string>
  <key>ProgramArguments</key>
  <array>
    <string>__REPO_ROOT__/.venv/bin/python</string>
    <string>-m</string>
    <string>src.watcher</string>
    <string>run</string>
    <string>__REPO_ROOT__/memory</string>
    <string>__REPO_ROOT__/AGENTS.md</string>
    <string>__REPO_ROOT__/SOUL.md</string>
    <string>__REPO_ROOT__/USER.md</string>
    <string>__REPO_ROOT__/IDENTITY.md</string>
    <string>__REPO_ROOT__/TOOLS.md</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>__HOME__</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/data/watcher.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/data/watcher.err.log</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
